# raspberry-pi/spark-device.service.template
[Unit]
Description=Spark device helper service (Wi-Fi/Bluetooth bridge for student page)
After=network.target bluetooth.target

[Service]
Type=simple
User=__SPARK_USER__
Group=__SPARK_USER__
ExecStart=/usr/bin/python3 /opt/spark/spark-device-service.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target